# hw6/monitor/monitor.service
[Unit]
Description=Key word monitor script

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/monitor
ExecStart=/usr/bin/monitor.sh
